# data/3rd-party-update.service
[Unit]
Description=Update 3rd-Party Software Content

[Service]
Type=oneshot
ExecStart=/usr/bin/swupd-3rd-party update
